# scheduled/crontab
# at 3am, run daily stuff
0 3 * * * bash /app/scheduled/daily.sh