# scheduled/crontab
# at 5am, run daily stuff
0 5 * * * bash /app/scheduled/daily.sh