# scheduled/crontab
*/1 * * * * cd /app && ./manage cron-tester